# rootfs/etc/cron.d/ssh-rotate-local-key
# Rotate local SSH key daily (SSH_ALLOW_LOCAL feature)
0 0 * * * root /usr/local/bin/ssh-rotate-local-key && cat /home/ubuntu/.ssh/id_ed25519.pub > /etc/ssh/authorized_keys
