# rootfs/etc/update-motd.d/80-openclaw
#!/bin/bash
# Skip if MOTD_SKIP is set (prevents recursion from openclaw wrapper)
[ -n "$MOTD_SKIP" ] && exit 0

source /etc/update-motd.d/lib.sh

echo -e "  🦞 ${BOLD}OpenClaw${NC}"
echo -e "  ${DIM}─────────────────────────────────────${NC}"

# Get health status - call directly if openclaw user, otherwise try sudo then su (non-login shells to avoid MOTD recursion)
OPENCLAW_BIN="/home/openclaw/.local/share/pnpm/openclaw"
if [ "$(id -un)" = "openclaw" ]; then
  HEALTH=$(timeout 5 "$OPENCLAW_BIN" health --json 2>/dev/null)
elif command -v sudo >/dev/null 2>&1; then
  HEALTH=$(timeout 5 sudo -n -u openclaw "$OPENCLAW_BIN" health --json 2>/dev/null)
else
  HEALTH=$(timeout 5 su openclaw -s /bin/bash -c "$OPENCLAW_BIN health --json" 2>/dev/null)
fi
if [ -n "$HEALTH" ]; then
  OK=$(echo "$HEALTH" | jq -r '.ok')
  if [ "$OK" = "true" ]; then
    printf "  ${CYAN}%-12s${NC} ${GREEN}✓ healthy${NC}\n" "Status:"
  else
    printf "  ${CYAN}%-12s${NC} ${RED}✗ unhealthy${NC}\n" "Status:"
  fi

  # Show configured channels
  CHANNELS=$(echo "$HEALTH" | jq -r '.channels | to_entries[] | select(.value.configured == true) | .key' 2>/dev/null | tr '\n' ', ' | sed 's/,$//')
  if [ -n "$CHANNELS" ]; then
    printf "  ${CYAN}%-12s${NC} ${YELLOW}%s${NC}\n" "Channels:" "$CHANNELS"
  else
    printf "  ${CYAN}%-12s${NC} ${DIM}none configured${NC}\n" "Channels:"
  fi

  # Show agent count
  AGENT_COUNT=$(echo "$HEALTH" | jq -r '.agents | length')
  printf "  ${CYAN}%-12s${NC} %s\n" "Agents:" "$AGENT_COUNT"
else
  printf "  ${CYAN}%-12s${NC} ${RED}✗ not running${NC}\n" "Status:"
fi

echo ""
